# home/dot_config/systemd/user/vibewatch.service.d/graphical-session.conf
# Drop-in: bind vibewatch to graphical-session lifecycle.
#
# Upstream ships WantedBy=default.target with only After=graphical-session.target.
# After= alone doesn't pull the target in, so vibewatch starts at user-session
# boot — before the compositor exists — crashes with "could not detect
# compositor", and burns through Restart=on-failure attempts. Adding PartOf=
# plus re-homing the enable symlink under graphical-session.target.wants/
# (see ../graphical-session.target.wants/symlink_vibewatch.service) starts it
# after Hyprland/Niri is up and stops it cleanly on session exit.
[Unit]
PartOf=graphical-session.target
After=graphical-session.target
