# home/dot_config/systemd/user/wayland-wm@.service.d/short-stop-timeout.conf
# Don't wait the full upstream 10s for the wayland-wm scope to drain on logout.
#
# Apps like Discord, Slack, and (sometimes) Chrome ignore SIGTERM and have to
# be SIGKILLed anyway — the 10s wait is dead time before SDDM can show the
# greeter. 3s gives well-behaved apps room to exit cleanly; everything else
# gets SIGKILLed sooner.
[Service]
TimeoutStopSec=3
