# home/dot_config/systemd/user/xdg-desktop-portal.service.d/graphical-session.conf
[Unit]
# Stock unit has WantedBy=default.target so the portal frontend outlives any
# one session. That leaves it alive during uwsm logout — and when it
# DBus-reactivates plasma-xdg-desktop-portal-kde for some last-second cleanup
# request, the KDE backend tries to connect to Wayland (already gone) and
# SIGABRTs (5× per logout in our case).
#
# Binding the frontend to graphical-session.target makes both stop together
# in dependency order, eliminating the re-activation race.
PartOf=graphical-session.target
After=graphical-session.target

[Service]
Slice=background-graphical.slice
