# home/dot_config/systemd/user/plasma-polkit-agent.service.d/graphical-session.conf
# Upstream sets PartOf=graphical-session.target but orders After=plasma-core.target;
# under uwsm we want it gated on the graphical session being reached instead.
# Upstream Slice=background.slice doesn't stop with the session — override.
[Unit]
After=graphical-session.target

[Service]
Slice=background-graphical.slice
