# home/dot_config/systemd/user/dropbox.service
[Unit]
Description=Dropbox
Documentation=https://help.dropbox.com/installs/linux
PartOf=graphical-session.target
After=graphical-session.target
Requisite=graphical-session.target

# Skip cleanly until the user has installed Dropbox via the bypass path.
ConditionPathExists=%h/.dropbox-dist/dropboxd

[Service]
Type=simple
ExecStart=%h/.dropbox-dist/dropboxd
Restart=on-failure
RestartSec=5

[Install]
WantedBy=graphical-session.target
