# home/dot_config/systemd/user/dropbox.service
[Unit]
Description=Dropbox
Documentation=https://help.dropbox.com/installs/linux
PartOf=graphical-session.target
After=graphical-session.target
Requisite=graphical-session.target

# Stay inactive (no crash-loop) until the user installs Dropbox via the
# bypass path — see commit 4a75437. The graphical-session.target.wants
# symlink can stay; systemd skips a unit whose condition is false.
ConditionPathExists=%h/.dropbox-dist/dropboxd

[Service]
Type=simple
ExecStart=%h/.dropbox-dist/dropboxd
Restart=on-failure
RestartSec=5

[Install]
WantedBy=graphical-session.target
